# emacs/rdm.service
[Unit]
Description=RTags daemon

Requires=rdm.socketp

[Service]
Type=simple
ExecStart=/home/keo9ren/.tools/scripts-and-more/emacs/rtags/bin/ -v --inactivity-timeout 300